# Script/git-auto-push/com.user.git-auto-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique label for this job -->
    <key>Label</key>
    <string>com.user.git-auto-push</string>
    
    <!-- Program to execute -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zyf/Desktop/Programmer/Script/git-auto-push/git-auto-push.sh</string>
    </array>
    
    <!-- Schedule: Check twice daily (every 12 hours) -->
    <!-- Combined with script's daily-once logic, only first execution will proceed -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/zyf/Library/Logs/git-auto-push.log</string>
    
    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/zyf/Library/Logs/git-auto-push.error.log</string>
    
    <!-- Run even when user is not logged in -->
    <key>RunAtLoad</key>
    <false/>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/zyf/Desktop/Programmer/Script</string>
    
    <!-- Keep job loaded even if it fails -->
    <key>KeepAlive</key>
    <false/>
    
    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
